# tasks/iswho/app/supervisord.conf
[supervisord]
nodaemon=true

[program:socat]
command=socat UNIX-LISTEN:/iswho/sock/iswho.sock,fork,mode=0666 TCP4-CONNECT:127.0.0.1:3000

[program:iswho]
command=java -jar /iswho/iswho.jar
